# crontab
12 15 * * sun /bin/systemctl restart device-farm.service

